# config/systemd/user/touchpad-screen-zoom.service
[Unit]
Description=Map touchpad scroll gestures to GNOME screen magnifier zoom
After=graphical-session.target
PartOf=graphical-session.target

[Service]
Type=simple
ExecStart=%h/.local/bin/touchpad-screen-zoom
Restart=on-failure
RestartSec=2

[Install]
WantedBy=default.target
