# config/autostart/copyq.desktop
[Desktop Entry]
Type=Application
Name=CopyQ
Exec=copyq
Icon=copyq
Comment=Clipboard manager
X-GNOME-Autostart-enabled=true
