# config/autostart/redshift-gtk.desktop
[Desktop Entry]
Type=Application
Name=Redshift
Exec=redshift-gtk
X-GNOME-Autostart-enabled=false
Hidden=true
